# eva.service
[Unit]
After=sys-subsystem-net-devices-eth0.device

[Service]
User=pi
PermissionsStartOnly=true
TimeoutStartSec=0
ExecStartPre=/bin/chmod o+rw /dev/vchiq
ExecStart=/home/pi/eva/eva -- -i 3 can0
KillMode=process
WorkingDirectory=/home/pi/eva

[Install]
WantedBy=multi-user.target
